# udev_rules/95-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/matt/.Xauthority", RUN+="/usr/bin/bash /usr/local/bin/monitor-test"
